# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:gunicorn]
command=/usr/local/bin/gunicorn bells.wsgi:application -b localhost:5000
directory=/code
